# deploy/com.morningstatic.prefrontal-coach.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  launchd agent that runs the native coaching tick (deploy/coach.sh):
  `prefrontal coach --deliver`. Fans over every enabled module (hyperfocus,
  outing/arrival, departure, self-care, task-paralysis, impulsivity,
  trip-tracking) plus the proactive overwhelm check, applies quiet-hours +
  debounce, and publishes what fires over ntfy/Pushover/TTS.

  This replaces four n8n poll workflows at once — coach-check, hyperfocus-check,
  departure-reminder, and panic-check. Before enabling it, DEACTIVATE those four
  in the n8n UI so a cue isn't sent twice (dedup/debounce catches most doubles,
  but don't rely on it).

  Install:
    cp deploy/com.morningstatic.prefrontal-coach.plist ~/Library/LaunchAgents/
    # then edit the paths below to match your username and install location
    launchctl load -w ~/Library/LaunchAgents/com.morningstatic.prefrontal-coach.plist
    # run the wrapper once by hand first to confirm it delivers:
    PREFRONTAL_HOME=/Users/tom/src/prefrontal deploy/coach.sh

  StartInterval fires every 60s (matching the old every-minute n8n poll). The
  tick self-suppresses via quiet hours + debounce, so a minute cadence only
  means "check often," not "nudge often."
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.morningstatic.prefrontal-coach</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/tom/src/prefrontal/deploy/coach.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/tom/src/prefrontal</string>

    <key>StartInterval</key>
    <integer>60</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/tom/Library/Logs/prefrontal.coach.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/tom/Library/Logs/prefrontal.coach.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <!-- Homebrew (Apple Silicon) on PATH so Ollama and other tools resolve. -->
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <!-- Install location; coach.sh cd's here so the adjacent .env loads. -->
        <key>PREFRONTAL_HOME</key>
        <string>/Users/tom/src/prefrontal</string>
        <!-- Handle to coach; leave unset to auto-pick when there's one user. -->
        <key>PREFRONTAL_USER</key>
        <string>Tom</string>
    </dict>
</dict>
</plist>
