# deploy/com.prefrontal.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  launchd agent that runs `prefrontal serve` always-on.

  Install (recommended): run `bash deploy/install-launchd.sh` — it fills the
  __PREFRONTAL_HOME__ / __HOME__ placeholders below with your real repo root and
  home directory, then loads the agent. The committed file is a template with no
  machine-specific paths.

  Manual alternative: substitute those placeholders yourself, then
    cp deploy/com.prefrontal.plist ~/Library/LaunchAgents/
    launchctl load -w ~/Library/LaunchAgents/com.prefrontal.plist

  WorkingDirectory must be the repo root so the `.env` next to it is loaded.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.prefrontal</string>

    <key>ProgramArguments</key>
    <array>
        <string>__PREFRONTAL_HOME__/.venv/bin/prefrontal</string>
        <string>serve</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__PREFRONTAL_HOME__</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <!-- Prefrontal holds one SQLite connection (db + -wal + -shm = 3 fds) per
         worker thread. The macOS default soft limit of 256 is far too low and
         the process can exhaust fds, 500ing DB endpoints with "unable to open
         database file". Raise it well clear of steady state. (The store also
         reaps dead worker threads' connections, so this is headroom, not a
         band-aid for an unbounded leak.) -->
    <key>SoftResourceLimits</key>
    <dict>
        <key>NumberOfFiles</key>
        <integer>4096</integer>
    </dict>
    <key>HardResourceLimits</key>
    <dict>
        <key>NumberOfFiles</key>
        <integer>8192</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/prefrontal.out.log</string>

    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/prefrontal.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <!-- Homebrew (Apple Silicon) on PATH so child processes resolve as expected. -->
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
