# deploy/com.prefrontal-mail.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  launchd agent that fetches + triages mail on a schedule (the no-n8n path).

  Unlike the always-on `prefrontal serve` agent, this is a *periodic* job: it
  runs `deploy/mail-fetch.sh` every StartInterval seconds, which calls
  `prefrontal mail fetch --all-users` — one job fans over every active user,
  fetching each user's own accounts from the per-user `sources` registry
  (`prefrontal mail add-source`). KeepAlive is intentionally absent — the job
  should exit after each run, not be restarted. launchd will not start a second
  copy while one is still running, so a slow IMAP pull never overlaps itself.

  Install (recommended): run `bash deploy/install-launchd.sh` — it fills the
  __PREFRONTAL_HOME__ / __HOME__ placeholders below with your real repo root and
  home directory, then loads the agent. The committed file is a template with no
  machine-specific paths.

  Manual alternative: substitute those placeholders yourself, then
    cp deploy/com.prefrontal-mail.plist ~/Library/LaunchAgents/
    launchctl load -w ~/Library/LaunchAgents/com.prefrontal-mail.plist

  WorkingDirectory must be the repo root so the `.env` next to it is loaded
  (PREFRONTAL_SECRET_KEY opens the sealed per-user IMAP passwords).

  PREREQUISITE: each user must have connected a mailbox first
  (`prefrontal mail --user <handle> add-source …`, or `import-env-sources` to
  migrate the legacy global MAIL_IMAP_* env). A user with no source is skipped —
  the global env mailbox is never fetched into another user's scope. Do a
  ONE-TIME heuristic backlog clear per account so the first run doesn't push a
  full inbox through the model:
    cd __PREFRONTAL_HOME__
    .venv/bin/prefrontal mail --user <handle> fetch --account personal --heuristic
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.prefrontal-mail</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__PREFRONTAL_HOME__/deploy/mail-fetch.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__PREFRONTAL_HOME__</string>

    <!-- Every 15 minutes. launchd skips runs while the Mac is asleep and fires
         once on wake if a scheduled interval elapsed. -->
    <key>StartInterval</key>
    <integer>900</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/prefrontal-mail.out.log</string>

    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/prefrontal-mail.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <!-- Homebrew (Apple Silicon) on PATH so child processes resolve as expected. -->
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
