# deploy/com.prefrontal-coach.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  launchd agent that runs the native coaching tick (deploy/coach.sh):
  `prefrontal coach --deliver`. Fans over every enabled module (hyperfocus,
  outing/arrival, departure, self-care, task-paralysis, impulsivity,
  trip-tracking) plus the proactive overwhelm check, applies quiet-hours +
  debounce, and publishes what fires over ntfy/Pushover/TTS.

  This replaces four n8n poll workflows at once — coach-check, hyperfocus-check,
  departure-reminder, and panic-check. Before enabling it, DEACTIVATE those four
  in the n8n UI so a cue isn't sent twice (dedup/debounce catches most doubles,
  but don't rely on it).

  One job covers everyone: coach.sh runs `prefrontal coach --deliver --all-users`,
  so there's no per-user handle to set. Delivery is per-user (a user without their
  own ntfy/Pushover target is never delivered to the operator's device).

  Install (recommended): run `bash deploy/install-launchd.sh` — it fills the
  __PREFRONTAL_HOME__ / __HOME__ placeholders below with your real repo root and
  home directory, then loads the agent. The committed file is a template with no
  machine-specific paths.

  Manual alternative: substitute those placeholders yourself, then
    cp deploy/com.prefrontal-coach.plist ~/Library/LaunchAgents/
    launchctl load -w ~/Library/LaunchAgents/com.prefrontal-coach.plist
    # run the wrapper once by hand first to confirm it delivers:
    deploy/coach.sh

  StartInterval fires every 60s (matching the old every-minute n8n poll). The
  tick self-suppresses via quiet hours + debounce, so a minute cadence only
  means "check often," not "nudge often."
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.prefrontal-coach</string>

    <key>ProgramArguments</key>
    <array>
        <string>__PREFRONTAL_HOME__/deploy/coach.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__PREFRONTAL_HOME__</string>

    <key>StartInterval</key>
    <integer>60</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/prefrontal.coach.log</string>

    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/prefrontal.coach.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <!-- Homebrew (Apple Silicon) on PATH so Ollama and other tools resolve. -->
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <!-- Install location; coach.sh cd's here so the adjacent .env loads. -->
        <key>PREFRONTAL_HOME</key>
        <string>__PREFRONTAL_HOME__</string>
    </dict>
</dict>
</plist>
